# files/systemd/grayhaven-ansible-runner.timer
[Unit]
Description=Run Grayhaven Systems LLC Ansible full playbook daily

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=15m
Persistent=true
Unit=grayhaven-ansible-runner.service

[Install]
WantedBy=timers.target
