# files/systemd/grayhaven-ansible-runner.service
[Unit]
Description=Grayhaven Systems LLC Ansible full playbook runner
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ansible
Group=ansible
WorkingDirectory=/home/ansible/ansible-runner
ExecStart=/usr/local/sbin/grayhaven-ansible-runner
# Keep convergence comfortably below the Grafana CPU alert threshold.
CPUAccounting=yes
CPUQuota=65%
