# roles/deploy_websites/templates/grayhaven-website-deploy-webhook.service.j2
[Unit]
Description=Grayhaven Systems LLC website deployment webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/grayhaven-website-deploy serve --config {{ deploy_websites_config_dir }}/deployments.json --host {{ deploy_websites_webhook_bind_host }} --port {{ deploy_websites_webhook_bind_port }} --runtime-dir {{ deploy_websites_runtime_dir }}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
