# roles/timetracker/templates/grayhaven-timetracker.container.j2
# Managed by Ansible. Local changes will be overwritten.

[Unit]
Description=Grayhaven Systems LLC Time Tracker
Wants=network-online.target
After=network-online.target

[Container]
Image={{ timetracker_image_repository }}@{{ timetracker.image_digest }}
ContainerName={{ timetracker_container_name }}
Pull=never
User={{ timetracker_uid }}:{{ timetracker_gid }}
PublishPort=127.0.0.1:{{ timetracker_host_port }}:{{ timetracker_container_port }}
Volume={{ timetracker.data_dir }}:/app/data:Z
Volume={{ timetracker.backup_dir }}:/app/backups:Z
Volume={{ timetracker.branding_dir }}:/app/branding:ro,Z
Volume={{ timetracker.secrets_dir }}:/run/secrets:ro,Z
Environment=BRANDING_PATH=/app/branding
Environment=BOOTSTRAP_USERS_FILE=/run/secrets/bootstrap_users
Environment=DATABASE_PATH=/app/data/{{ timetracker_database_filename }}
Environment=PUBLIC_BASE_URL=https://{{ timetracker.hostname }}
Environment=SECRET_KEY_FILE=/run/secrets/flask_secret_key
Environment=SESSION_COOKIE_SECURE=true
Environment=SQLCIPHER_PASSPHRASE_FILE=/run/secrets/sqlcipher_passphrase
Environment=TRUSTED_HOSTS={{ timetracker.hostname }}
Environment=TRUSTED_PROXY_COUNT=1
Environment=TZ={{ grayhaven_timezone }}
DropCapability=all
NoNewPrivileges=true
PidsLimit=128
ReadOnly=true
Tmpfs=/tmp:rw,noexec,nosuid,nodev,size=16m
HealthCmd=python -c "import urllib.request; urllib.request.urlopen('http://127.0.0.1:8000/health', timeout=3)"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s
Notify=healthy

[Service]
Restart=always
RestartSec=5s
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
